# install/etc/supervisor/conf.d/dokku-api.conf
[program:dokku-api]
command=/home/vagrant/.golang/bin/dokku-api
autostart=true
autorestart=unexpected
stderr_logfile=/var/log/dokku-api.err.log
stdout_logfile=/var/log/dokku-api.out.log
